# deploy/icom-demod.service
[Unit]
Description=IC-705 IF IQ demodulator (pw-record -> GNU Radio -> TCP :7355)
# pipewire живёт в user-сессии sinty (linger включён) — ждём её
After=user@1000.service sound.target
Requires=user@1000.service

[Service]
User=sinty
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/bin/bash /home/sinty/run_pipeline.sh
Restart=always
RestartSec=3
Nice=-10

[Install]
WantedBy=multi-user.target
